# python/docker/docker_com/django-uwsgi-nginx/supervisor-app.conf
[program:app-uwsgi]

command = '/usr/local/bin/uwsgi --ini /code/com/com.ini'


[program:nginx-app]
command = '/usr/sbin/nginx'
